# pi_vision/pi-sonar.service
# systemd unit: auto-start the HC-SR04 distance sender on boot.
#
# Runs sonar_main.py once the network is up, so the cane reports distance as
# soon as it powers on — no laptop / SSH needed. sonar_main.py auto-detects the
# phone as the default gateway (on the phone's hotspot) and reconnects on its
# own if the app isn't open yet, so the Pi can boot first and just wait.
#
# This is the WiFi replacement for the ESP32 ultrasonic path; it runs alongside
# pi-vision.service (camera) — different GPIOs, different TCP port.
#
# --- INSTALL (run on the Pi, once) ---
#   1. Edit User= and WorkingDirectory=/ExecStart= below to match YOUR Pi.
#   2. GPIO library (the reader uses lgpio — native on Bookworm, no daemon;
#      pigpio is NOT used and is no longer in Bookworm's repos):
#        sudo apt install -y python3-lgpio
#   3. sudo cp pi_vision/pi-sonar.service /etc/systemd/system/pi-sonar.service
#   4. sudo systemctl daemon-reload
#   5. sudo systemctl enable --now pi-sonar.service   # start now + on every boot
#   6. systemctl status pi-sonar.service              # check it's running
#      journalctl -u pi-sonar.service -f              # follow its logs
#
#   To stop / disable:  sudo systemctl disable --now pi-sonar.service

[Unit]
Description=Smart Cane - Pi Zero HC-SR04 distance sender
# Wait until the network is up (joined the hotspot) before starting. The reader
# uses lgpio directly (no daemon), so there's nothing else to order against.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# CHANGE THESE to match your Pi:
#   User    = the login you set up in Raspberry Pi Imager (NOT necessarily "pi").
#             Find it with: whoami
#   WorkingDirectory / ExecStart = the absolute path to this pi_vision folder.
#             Find it with: cd pi_vision && pwd
User=apn7
WorkingDirectory=/home/apn7/Thesis_Pi_Zero/pi_vision
ExecStart=/usr/bin/python3 /home/apn7/Thesis_Pi_Zero/pi_vision/sonar_main.py

# Assistive device: the sender must ALWAYS come back, even after a clean-but-
# unexpected exit. `systemctl stop` still stops it (manual stops don't restart).
Restart=always
RestartSec=3
# Never give up: the default rate limit (5 starts/10 s) would permanently fail
# the unit during a long "waiting for hotspot" boot loop.
StartLimitIntervalSec=0

# Send Python logs straight to the journal (journalctl) without buffering.
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
