# pi_vision/pi-sonar.service
# systemd unit: auto-start the cane's distance loop on boot.
#
# Runs sonar_main.py, which now does two jobs:
#   1. STANDALONE alerts — classifies every HC-SR04 reading on the Pi (same
#      thresholds + hysteresis as the app) and drives the buzzer + vibration
#      motor. Works with no phone, no app, no network.
#   2. Phone telemetry — a background thread auto-detects the phone (hotspot
#      gateway or AP client) and streams the raw readings; the app classifies
#      them with identical rules, so phone and cane alert in lockstep.
#
# This is the WiFi replacement for the ESP32 ultrasonic path; it runs alongside
# pi-vision.service (camera) — different GPIOs, different TCP port.
#
# --- INSTALL (run on the Pi, once) ---
#   1. Edit User= and WorkingDirectory=/ExecStart= below to match YOUR Pi.
#   2. GPIO library (the reader uses lgpio — native on Bookworm, no daemon;
#      pigpio is NOT used and is no longer in Bookworm's repos):
#        sudo apt install -y python3-lgpio
#   3. sudo cp pi_vision/pi-sonar.service /etc/systemd/system/pi-sonar.service
#   4. sudo systemctl daemon-reload
#   5. sudo systemctl enable --now pi-sonar.service   # start now + on every boot
#   6. systemctl status pi-sonar.service              # check it's running
#      journalctl -u pi-sonar.service -f              # follow its logs
#
#   To stop / disable:  sudo systemctl disable --now pi-sonar.service

[Unit]
Description=Smart Cane - Pi Zero distance alerts + phone sender
# Deliberately NO network-online dependency: the buzzer/vibration reflex must
# come up even if WiFi never does (standalone cane). The network half of
# sonar_main.py polls for the phone forever on its own thread, so starting
# before the network is up is the designed-for case, not a race.

[Service]
Type=simple

# CHANGE THESE to match your Pi:
#   User    = the login you set up in Raspberry Pi Imager (NOT necessarily "pi").
#             Find it with: whoami
#   WorkingDirectory / ExecStart = the absolute path to this pi_vision folder.
#             Find it with: cd pi_vision && pwd
User=apn7
WorkingDirectory=/home/apn7/Thesis_Pi_Zero/pi_vision
ExecStart=/usr/bin/python3 /home/apn7/Thesis_Pi_Zero/pi_vision/sonar_main.py

# Assistive device: the sender must ALWAYS come back, even after a clean-but-
# unexpected exit. `systemctl stop` still stops it (manual stops don't restart).
Restart=always
RestartSec=3
# Never give up: the default rate limit (5 starts/10 s) would permanently fail
# the unit during a long "waiting for hotspot" boot loop.
StartLimitIntervalSec=0

# Send Python logs straight to the journal (journalctl) without buffering.
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
